# systemd/sdmnh-orchestrator.timer
[Unit]
Description=Daily SDMNH render trigger

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
