# systemd/sdmnh-morning.timer
[Unit]
Description=SDMNH Morning Render (11:00 UTC → upload ~2pm UTC)
Requires=sdmnh-orchestrator@morning.service

[Timer]
OnCalendar=*-*-* 11:00:00 UTC
Persistent=true
Unit=sdmnh-orchestrator@morning.service

[Install]
WantedBy=timers.target
